# scripts/crontab/crontab.1
* * * * * cd /data/code/i_service/scripts && php orderRemind.php >> ./logs/orderRemind.log 2>&1 &
