# scripts/crontab/crontab.1
* * * * * cd /data/code/i_service/scripts && /usr/local/bin/php orderRemind.php >> ../logs/orderRemind.log 2>&1 &
